# etc/systemd/system/paste_capsule.socket
[Unit]
Description=gunicorn socket for paste_capsule

[Socket]
ListenStream=/run/paste_capsule.sock

[Install]
WantedBy=sockets.target
